# 50-winpath.rules
ACTION=="add|remove", KERNEL=="sd*", SUBSYSTEMS=="scsi", RUN+="sudo -H -u otheruser bash -c '/usr/bin/winpath'"
