# 50-winpath.rules
ACTION=="add|remove", KERNEL=="sd*", SUBSYSTEMS=="scsi", RUN+="runuser -l  winpath -c '/usr/bin/winpath'"
