# 50-winpath.rules
ACTION=="add|remove", KERNEL=="sd*", SUBSYSTEMS=="scsi", RUN+="/usr/bin/winpath"
