# 11.10/winpath.rules
SUBSYSTEM=="block", KERNEL=="s*[0-9]", RUN+="/usr/bin/winpath"
